# userrtm-gateway.service
[Unit]
Description=UserrTM Gateway
After=network-online.target nginx.service
Wants=network-online.target
Requires=nginx.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/userrtm-gateway
EnvironmentFile=-/etc/userrtm-gateway/userrtm-gateway.env
ExecStart=/opt/userrtm-gateway/userrtm-gateway
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
